# deploy/systemd/115-sync-downloader.service
[Unit]
Description=115 Sync Console differential downloader
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=media
Group=media
WorkingDirectory=/opt/115-sync-console
EnvironmentFile=/etc/115-sync-console.env
ExecStart=/opt/115-sync-console/.venv/bin/python /opt/115-sync-console/src/incremental_download.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
