# deploy/systemd/115-sync-dashboard.service
[Unit]
Description=115 Sync Console dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=media
Group=media
WorkingDirectory=/opt/115-sync-console
EnvironmentFile=/etc/115-sync-console.env
ExecStart=/opt/115-sync-console/.venv/bin/python /opt/115-sync-console/src/server.py
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
